# src/Utility/multichain_node_setup/multichain.service
[Unit]
Description=Multichain service

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/local/bin/multichaind MainChain

[Install]
WantedBy=multi-user.target